# systemdFiles/pacupdates.timer
[Unit]
Description=Pacupdates timer

[Timer]
OnCalendar=*:0/15
OnBootSec=10s

[Install]
WantedBy=timers.target